# eq_telescope/files/pigpiod.service
[Unit]
Description=Pigpiod daemon 

[Service]
ExecStart=/usr/bin/pigpiod -g
Restart=always

[Install]
WantedBy=multi-user.target
